# .docker/supervisor.d/horizon.conf
[program:laravel-horizon]
process_name=%(program_name)s
directory=/var/www/html
command=php artisan horizon
autostart=true
autorestart=true
stderr_logfile=/var/www/html/logs/%(program_name)s.err.log
stdout_logfile=/var/www/html/logs/%(program_name)s.out.log